# rootfs/etc/s6/minio/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

exec /usr/bin/minio server /var/lib/minio
